# crontab
# Metadata update
* */12 * * * /usr/local/bin/php /app/cli/sintoniza  2>&1 | tee -a /tmp/sintoniza.log